# ansible/roles/shipyard/templates/etc/profile.d/vault.sh.j2
#!/bin/bash

export VAULT_ADDR=http://{{ vault_host }}:{{ vault_port }}
